# luci-app-iwand/root/etc/init.d/iwand
#!/bin/sh /etc/rc.common
# iwand - Panabit iWAN SD-WAN client (UCI-aware init script)

START=99
STOP=10
USE_PROCD=1

PROG=/usr/sbin/iwand
CONF_DIR=/var/run/iwand
CONF_FILE="$CONF_DIR/iwan.conf"

# Reject values containing newlines, brackets, or control characters
validate_value() {
	case "$1" in
		*[[:cntrl:]]* | *'['* | *']'*) return 1 ;;
	esac
	return 0
}

generate_config() {
	local enabled interface server username password port mtu encrypt
	local pipeid pipeidx srlinks srpassword srencryptmode

	uci_validate_section iwand iwand "main" \
		'enabled:bool:0' \
		'interface:string:iwan0' \
		'server:string' \
		'username:string' \
		'password:string' \
		'port:port:4567' \
		'mtu:range(46,1600):1400' \
		'encrypt:bool:0' \
		'pipeid:uinteger:0' \
		'pipeidx:uinteger:0' \
		'srlinks:string' \
		'srpassword:string' \
		'srencryptmode:range(0,2):0' || {
		echo "validation failed" >&2
		return 2
	}

	[ "$enabled" = "1" ] || return 1
	[ -n "$server" ] || { echo "server not set" >&2; return 2; }
	[ -n "$username" ] || { echo "username not set" >&2; return 2; }
	[ -n "$password" ] || { echo "password not set" >&2; return 2; }

	# Sanitize free-form values
	validate_value "$interface" || { echo "invalid interface name" >&2; return 2; }
	validate_value "$server" || { echo "invalid server value" >&2; return 2; }
	validate_value "$username" || { echo "invalid username" >&2; return 2; }
	validate_value "$password" || { echo "invalid password" >&2; return 2; }

	install -d -m 0700 "$CONF_DIR"

	# Write config with restricted permissions (printf avoids shell expansion in values)
	(
		umask 077
		printf '[%s]\nserver=%s\nusername=%s\npassword=%s\nport=%s\nmtu=%s\nencrypt=%s\n' \
			"$interface" "$server" "$username" "$password" "$port" "$mtu" "$encrypt" \
			> "$CONF_FILE"
	)

	[ "$pipeid" -gt 0 ] 2>/dev/null && echo "pipeid=$pipeid" >> "$CONF_FILE"
	[ "$pipeidx" -gt 0 ] 2>/dev/null && echo "pipeidx=$pipeidx" >> "$CONF_FILE"
	[ -n "$srlinks" ] && validate_value "$srlinks" && echo "srlinks=$srlinks" >> "$CONF_FILE"
	[ -n "$srpassword" ] && validate_value "$srpassword" && echo "srpassword=$srpassword" >> "$CONF_FILE"
	[ "$srencryptmode" -gt 0 ] 2>/dev/null && echo "srencryptmode=$srencryptmode" >> "$CONF_FILE"

	chmod 0600 "$CONF_FILE"
	return 0
}

start_service() {
	generate_config
	local ret=$?
	# ret=1 means disabled (not an error), ret=2 means config error
	[ "$ret" = "1" ] && return 0
	[ "$ret" != "0" ] && return 1

	[ -e /dev/net/tun ] || {
		mkdir -p /dev/net
		mknod /dev/net/tun c 10 200
	}
	lsmod | grep -q '^tun ' || modprobe tun 2>/dev/null

	procd_open_instance
	procd_set_param command "$PROG" -F -f "$CONF_FILE"
	procd_set_param respawn 3600 5 5
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_close_instance
}

stop_service() {
	rm -f "$CONF_FILE"
}

service_triggers() {
	procd_add_reload_trigger "iwand"
	procd_add_validation validate_config
}

validate_config() {
	uci_validate_section iwand iwand "main" \
		'enabled:bool:0' \
		'interface:string:iwan0' \
		'server:string' \
		'username:string' \
		'password:string' \
		'port:port:4567' \
		'mtu:range(46,1600):1400' \
		'encrypt:bool:0' \
		'pipeid:uinteger:0' \
		'pipeidx:uinteger:0' \
		'srlinks:string' \
		'srpassword:string' \
		'srencryptmode:range(0,2):0'
}

reload_service() {
	stop
	start
}
